# bin/service
#!/bin/bash
#
# chkconfig: 35 90 12
# description: DOM0 server
#
# Get function from functions library
#. /etc/init.d/functions
# Start the service DOM0
DIR="$( cd "$( dirname "${BASH_SOURCE[0]}" )" && pwd )"
start() {
        echo -n "DOM0 server started."
        $DIR/memory_server &
        ### Create the lock file ###
        touch /var/lock/TEST
        #success $"TEST server startup"
        echo
}
# Restart the service TEST
stop() {
        pkill memory_server
        rm -f /var/lock/TEST
        echo "DOM0 server stopped."
}
### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status FOO
        ;;
  restart|reload|condrestart)
        stop
        start
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        exit 1
esac
exit 0
